# packagesNotesAndTroubleShootings/StartupServiceTemplate.service
[Unit]
After=network.service

[Service]
ExecStart=/opt/StartupScript/GuestMountVHDX.sh

[Install]
WantedBy=default.target